# roles/cliente_iscsi_linux/files/mnt.automount
Description=Automontaje de Volumen de iscsi

[Automount]
Where=/mnt

[Install]
WantedBy=multi-user.target